# bot.conf
[program:mydiplom1]
directory = /home/ubuntu/MyDiplom
command =python3 app.py

autostart = true
autorestart = true

environment = HOME = "/home/ubuntu",USER = "ubuntu"

stderr_logfile =/home/ubuntu/MyDiplom/logfile_err.log
stdout_logfile =/home/ubuntu/MyDiplom/logfile.log